# images/ruby/etc/profile.d/tf.ruby.sh
export PATH="./.bundle/bin:$HOME/.gem/ruby/2.0.0/bin:$PATH"
