# images/go/etc/profile.d/tf.go.sh
export GOPATH=$HOME/workspace
export PATH=$GOPATH/bin:$PATH
